# deploy/eggday-web.service
[Unit]
Description=Egg Day competition website
After=network-online.target
Wants=network-online.target

[Service]
User=harsh
WorkingDirectory=/srv/eggday/egg-day-worker
# Single worker on purpose: the TTL cache is in-process.
ExecStart=/home/egger/.local/bin/uv run uvicorn web.app:app --host 127.0.0.1 --port 8100 --workers 1
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
